# mnt-exthd.automount
[Unit]
Description=Automount usb

[Automount]
Where=/mnt/exthd

[Install]
WantedBy=multi-user.target
